# examples/backup_servers.cron
#/etc/cron.d/backup_servers
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
40 01 * * * root LOG=/var/log/backup_servers.log; backup_servers --conf /etc/backup_servers.conf > $LOG; backup_servers --log2err $LOG
