# examples/backup_servers.cron
40 02 * * * root LOG=/var/log/backup_servers.log (cd /usr/local/src/backup_servers && git pull) > $LOG; backup_servers --du --conf /home/backup/backup.conf >> $LOG; backup_servers --log2err $LOG
